# system/laravelTask.timer
[Unit]
Description=laravel Task

[Timer]
OnCalendar=*:0/1
Unit=laravelTask.service

[Install]
WantedBy=timers.target
